# backup/config/systemd/kiro-auth-refresh.timer
[Unit]
Description=Hourly Kiro auth token refresh
Documentation=https://github.com/alexandeism/kirolink
RefuseManualStart=no

[Timer]
OnCalendar=hourly
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
