# backup/config/systemd/kiro-gateway.service
[Unit]
Description=Kiro Gateway - AI Proxy for Amazon Q Developer
After=network.target
BindsTo=kiro-tokend.service
After=kiro-tokend.service

[Service]
Type=simple
WorkingDirectory=/home/x1/Documents/proxy/kiro-gateway
Environment=PROXY_API_KEY=kiro-gateway-8333
Environment=KIRO_CLI_DB_FILE=~/.local/share/kiro-cli/data.sqlite3
Environment=PATH=/home/x1/Documents/proxy/kirolink:/home/x1/.local/bin:/home/x1/.npm-global/bin:/home/x1/go/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/x1/Documents/proxy/kiro-gateway/.venv/bin/python main.py --port 8333
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
